# conf/gunicorn.conf
# supervisord config file for encyc-rg

[program:encycrg]
user=encyc
directory=/opt/encyc-rg/encycrg
command=/opt/encyc-rg/venv/encycrg/bin/gunicorn encycrg.wsgi:application -w 5 -b 0.0.0.0:8081
autostart=true
autorestart=true
redirect_stderr=True
